# launchd/com.trs0817.claude-code-telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.trs0817.claude-code-telegram</string>

    <!--
        Runs the env-sourcing launcher, which reads /etc/claude-code-telegram.env
        and execs the bot. launchd has no EnvironmentFile equivalent, so a wrapper
        script handles credential injection without exposing secrets in the plist.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>__REPLACE_LAUNCHER__</string>
    </array>

    <!-- Run as the configured user (never root) -->
    <key>UserName</key>
    <string>__REPLACE_USER__</string>

    <!-- Start on load; keep alive after unexpected exits -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!--
        Throttle restart attempts to avoid rapid crash loops.
        Mirrors systemd RestartSec=10.
    -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Log to a persistent file. View with: tail -f /var/log/claude-code-telegram.log -->
    <key>StandardOutPath</key>
    <string>/var/log/claude-code-telegram.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/claude-code-telegram.log</string>

</dict>
</plist>
